# backend/service/muri-app.service
[Unit]
Description=muri app python service

[Service]
WorkingDirectory=/home/muri-app
ExecStart=/usr/bin/python3 /home/muri-app/backend/muri_app_main.py
Restart=always
StandardOutput=file:/home/muri-app/std.out
StandardError=file:/home/muri-app/err.out

[Install]
WantedBy=multi-user.target
